# user/fvwm.target
[Unit]
Description=FVWM
Wants=xorg.target
Wants=fvwm.service
After=fvwm.service
AllowIsolate=true

[Install]
WantedBy=default.target
