# user/fvwm.target
[Unit]
Description=FVWM
Wants=fvwm.service
After=fvwm.service
AllowIsolate=true

[Install]
WantedBy=default.target
