# user/gnome-keyring.service
[Unit]
Description=GNOME Keyring Daemon
IgnoreOnIsolate=true

[Service]
ExecStart=/usr/bin/gnome-keyring-daemon --start -f
