# user/notification-daemon.service
[Unit]
Description=Notification Daemon

[Service]
ExecStart=/usr/libexec/notification-daemon

[Install]
WantedBy=default.target
